# deploy/gunicorn.service
[Unit]
Description=Food Delivery Gunicorn
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/food-delivery/backend
EnvironmentFile=/var/www/food-delivery/.env
ExecStart=/var/www/food-delivery/venv/bin/gunicorn food_delivery.wsgi:application --bind 127.0.0.1:8000 --workers 3
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
